# crontab/checkAccTime.sh
#!/bin/sh


php -f ~/checkAccRemove.php > results.txt

